# crontab
*/10 * * * * zsh ~/Code/dotfiles/bin/backup.sh
