# crontab
#min hour day  mon  wday
0    */5  *    *    *    zsh ~/Code/dotfiles/bin/backup
0    */5  */5  *    *    zsh ~/Code/dotfiles/bin/rmbackups
*/5  *    *    *    *    zsh ~/Code/dotfiles/bin/sync_mail
